# roles/desktop/templates/mbsync/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
User=%I
Group=%I
ExecStart=/usr/bin/mbsync -Va
