# roles/desktop/templates/mbsync/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
User=%i
Group=%i
ExecStart=/usr/bin/mbsync -Va
