# templates/home/user/.config/systemd/user/offlineimap.timer
[Timer]
OnUnitInactiveSec=300s
Unit=offlineimap.service

[Install]
WantedBy=default.target
