# roles/mail/templates/.config/systemd/user/mbsync.timer
[Timer]
OnUnitInactiveSec=300s
Unit=mbsync.service

[Install]
WantedBy=default.target
